# system/supervisor/conf.d/ssh-agent.conf
[program:ssh-agent]
command=/bin/bash -c 'rm -f /tmp/ssh-agent.sock && exec ssh-agent -a /tmp/ssh-agent.sock -d'
user=root
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/ssh-agent.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/ssh-agent.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1

[program:ssh-agent-keys]
command=/usr/local/bin/ssh-agent-add-keys.sh
user=root
autostart=true
autorestart=false
startsecs=0
priority=20
stdout_logfile=/var/log/supervisor/ssh-agent-keys.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/ssh-agent-keys.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
