# system/supervisor/conf.d/tailscale-up.conf
[program:tailscale-up]
command=/usr/local/bin/tailscale-up.sh
user=root
autostart=true
autorestart=false
startsecs=0
priority=20
stdout_logfile=/var/log/supervisor/tailscale-up.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/tailscale-up.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
