# containers/sandbox/supervisor/conf.d/archon-serve.conf
[program:archon-serve]
command=archon serve --host 0.0.0.0 --port 3090
directory=/home/gem/projects
environment=HOME="/home/gem",CLAUDE_BIN_PATH="/usr/bin/claude",TMPDIR="/home/gem/.cache/archon-tmp",PATH="/usr/local/cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
user=gem
autostart=true
autorestart=true
priority=50
stdout_logfile=/tmp/archon-serve.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/tmp/archon-serve.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
